# deploy/crypto-mm.service
[Unit]
Description=Crypto Trading Desk Intern
After=network.target
Wants=network-online.target

[Service]
User=ubuntu
WorkingDirectory=/opt/crypto-trading-desk-intern
EnvironmentFile=/opt/crypto-trading-desk-intern/.env
ExecStart=/opt/crypto-trading-desk-intern/.venv/bin/python -m crypto_mm.main
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

